# packages/squeezelite/src/squeezelite.service
[Unit]
Description=Squeezelite
Wants=network.target
After=network.target sound.target pipwire-server.service
Before=initial-volume.service

[Service]
Type=simple
Environment=HOME=/var/lib/squeezelite
User=squeezelite
Group=audio
ExecStartPre=-/usr/bin/config-volume --store
ExecStart=/usr/bin/start-squeezelite
ExecStartPost=sleep 2
ExecStartPost=-/usr/bin/config-volume --restore
StandardOutput=journal
Restart=always
RestartSec=5
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
