# buildroot/package/dsptoolkit/sigmatcp.service
[Unit]
Description=SigmaTCP Server for HiFiBerry DSP
Wants=network.target sound.target
After=network.target sound.target

[Service]
Type=notify
ExecStartPre=/opt/hifiberry/bin/bootmsg "Starting DSP server"
ExecStartPre=/opt/hifiberry/bin/create-dspvolume
ExecStart=/usr/bin/sigmatcpserver --alsa  ${SIGMATCPSERVER_ARGS}
ExecStartPost=/bin/sleep 3
StandardOutput=journal
Restart=always
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
