# buildroot/package/hifiberry-tools/configure-players.service
[Unit]
Description=Reconfigure players
After=sound.target 

[Service]
Type=Oneshot
ExecStart=/opt/hifiberry/bin/reconfigure-players --reboot
StandardOutput=journal

[Install]
WantedBy=multi-user.target
