# buildroot/package/hifiberry-tools/initial-volume.service
[Unit]
Description=Set initial volume
Wants=sound.target
After=sound.target squeezelite.service

[Service]
Type=oneshot
Environment=HOME=/root
ExecStartPre=/opt/hifiberry/bin/bootmsg "Setting initial volume"
ExecStart=/opt/hifiberry/bin/set-initial-volume
StandardOutput=journal

[Install]
WantedBy=multi-user.target
