# buildroot/package/hifiberry-tools/internet.service
# This will detect if the system is connected to the internet independent from the "network-online" target
[Unit]
Description=Check Internet Service Availability
After=network.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=/opt/hifiberry/bin/check-internet.py https://www.cloudflare.com https://www.google.com https://www.microsoft.com --timeout 60
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
