# buildroot/package/hifiberry-tools/store-volume.service
[Unit]
Description=Store volume
Wants=sound.target
After=sound.target

[Service]
Type=oneshot
ExecStart=/opt/hifiberry/bin/store-volume

