# buildroot/package/hifiberry-squeezelite/squeezelite.service
[Unit]
Description=Squeezelite
Wants=network.target
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/squeezelite -o softvol -M HiFiBerry -N /var/squeezelite/squeezelite.name
StandardOutput=journal
Restart=Always
RestartSec=5

[Install]
WantedBy=multi-user.target
