# buildroot/package/hifiberry-ympd/ympd.service
[Unit]
Description=ympd MPD web interface
Requires=mpd.service
After=mpd.service

[Service]
Type=simple
ExecStartPre=/opt/hifiberry/bin/bootmsg "Starting MPD web interface"
ExecStart=/usr/bin/ympd -w 9000
StandardOutput=journal
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
