# packages/TODO/hifiberry-tools/services/configure-players.service
[Unit]
Description=Reconfigure players
After=sound.target hifiberry-detect.service local-fs.target sigmatcp.service
Before=shairport-sync.service squeezelite.service spotify.service mpd.service raat.service

[Service]
Type=oneshot
Environment=HOME=/root
ExecStartPre=/opt/hifiberry/bin/bootmsg "Configuring players"
ExecStart=/opt/hifiberry/bin/reconfigure-players
StandardOutput=journal

[Install]
WantedBy=multi-user.target
