# packages/TODO/hifiberry-tools/services/configure-system.service
[Unit]
Description=Reconfigure system
After=boot.mount local-fs.target
Before=network.target

[Service]
Type=oneshot
ExecStartPre=/opt/hifiberry/bin/bootmsg "Configuring system"
ExecStart=/opt/hifiberry/bin/configure-system
StandardOutput=journal

[Install]
WantedBy=multi-user.target
