# buildroot/package/hifiberry-automount/mount-data.service
[Unit]
Description=Mount /data
After=local-fs.target datacollector.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/opt/hifiberry/bin/mount-data.sh
