# packages/librespot/librespot-docker-build/librespot.service
[Unit]
Description=Librespot Spotify Connect Player
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=librespot
Group=audio
Environment=RUST_BACKTRACE=1
ExecStart=/usr/bin/librespot --name HiFiBerry --backend alsa --device default --bitrate 320 --disable-audio-cache
Restart=always
RestartSec=10
StartLimitInterval=30
StartLimitBurst=5

[Install]
WantedBy=multi-user.target