# buildroot/package/hifiberry-psplash/psplash-start.service
[Unit]
Description=Starts Psplash Boot screen
DefaultDependencies=no

[Service]
ExecStart=/usr/bin/psplash -n

[Install]
WantedBy=sysinit.target
