# packages/pipewire-server/postinstall
#!/bin/bash
set -e

# Create the runtime directory with the correct permissions
if [ ! -d /run/pipewire ]; then
    mkdir -p /run/pipewire
    chmod 1777 /run/pipewire
fi

# Create tmpfiles configuration for persistence
echo "d /run/pipewire 0777 root root -" > /etc/tmpfiles.d/pipewire.conf
systemd-tmpfiles --create /etc/tmpfiles.d/pipewire.conf

# Enable and start the PipeWire service
echo "Enabling and starting pipewire-server.service..."
systemctl daemon-reload
systemctl enable pipewire-server.service
systemctl start pipewire-server.service

# Remove PipeWire files and directories in /etc/systemd/user
echo "Removing PipeWire-related files and symlinks from /etc/systemd/user..."
if [ -d /etc/systemd/user ]; then
    find /etc/systemd/user -name "pipewire*.service" -exec rm -f {} \; || true
    find /etc/systemd/user -name "pipewire*.socket" -exec rm -f {} \; || true
    find /etc/systemd/user -name "pipewire*.target" -exec rm -f {} \; || true
    find /etc/systemd/user -type d -name "pipewire*.wants" -exec rm -rf {} \; || true
fi

echo "PipeWire-related files and directories removed from /etc/systemd/user."

rm -rf /run/user/1000/pipewire-0*

exit 0

