# packages/shairport-sync/src/debian/prerm
#!/bin/bash
set -e

# Stop and disable systemd services on package removal
if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    if command -v systemctl >/dev/null; then
        systemctl stop shairport.service || true
        systemctl disable shairport.service || true
        systemctl stop nqptp.service || true
        systemctl disable nqptp.service || true
    fi
fi

#DEBHELPER#
