# packages/shairport-sync/src/debian/prerm
#!/bin/bash
set -e

# Stop and disable systemd services on package removal
if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    if [ -d /run/systemd/system ]; then
        deb-systemd-invoke stop shairport.service >/dev/null || true
        deb-systemd-helper disable shairport.service >/dev/null || true
        deb-systemd-invoke stop nqptp.service >/dev/null || true
        deb-systemd-helper disable nqptp.service >/dev/null || true
    fi
fi

#DEBHELPER#
