# buildroot/package/spotifyd/spotify.service
[Unit]
Description=Spotify connect
Wants=network-online.target
After=network.target sound.target dbus.service

[Service]
Type=simple
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket
ExecStart=/usr/bin/spotifyd --no-daemon 
StandardOutput=journal
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
