# buildroot/package/wisa-scripts/wisa.service
[Unit]
Description=WISA init
ConditionPathExists=/opt/wisa/etc/room.cfg

[Service]
Type=oneshot
ExecStartPre=/opt/hifiberry/bin/bootmsg "Starting WiSA"
ExecStart=/opt/wisa/bin/system_startup.py

[Install]
WantedBy=multi-user.target
