# packages/spotifyd/spotifyd-docker-build/systemd/spotifyd.service
[Unit]
Description=HiFiBerry Spotify Daemon
Documentation=https://github.com/Spotifyd/spotifyd
After=network.target sound.target
Wants=network.target sound.target

[Service]
Type=simple
User=spotifyd
Group=audio
ExecStart=/usr/bin/spotifyd --no-daemon
Restart=on-failure
RestartSec=10

# Security hardening
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/var/run

[Install]
WantedBy=multi-user.target