# packages/pipewire/src/wireplumber-system.service
[Unit]
Description=WirePlumber Session Manager (System)
Documentation=man:wireplumber(1)
After=pipewire.service
Requires=pipewire.service

[Service]
Type=simple
ExecStart=/usr/bin/wireplumber
Restart=on-failure
User=pipewire
Group=audio
RuntimeDirectory=pipewire
RuntimeDirectoryMode=0755
Environment=PIPEWIRE_RUNTIME_DIR=/run/pipewire
CacheDirectory=wireplumber
StateDirectory=wireplumber
LogsDirectory=wireplumber
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target