# buildroot/package/hifiberry-updater/updater.timer
[Unit]
Description=Install upgrades

[Timer]
OnCalendar=03:00
Unit=updater.service

[Install]
WantedBy=timers.target  
