# buildroot/package/hifiberry-updater/updater.service
[Unit]
Description=Install updates
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/hifiberry/bin/update

