# buildroot/package/hifiberry-updater/updater.service
[Unit]
Description=Install updates
After=internet.service

[Service]
Type=oneshot
ExecStart=/opt/hifiberry/bin/update

