# buildroot/package/hifiberry-watchdog/watchdog.service
[Unit]
Description=HiFiBerry watchdog
After=network.target
ConditionPathExists=!/boot/no_watchdog

[Service]
Type=simple
ExecStartPre=/opt/hifiberry/bin/bootmsg "Starting watchdog"
WorkingDirectory=/opt/watchdog
ExecStart=/bin/python /opt/watchdog/watchdog.py
StandardOutput=journal
Restart=always
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
