# packages/raat/src/raat.service
[Unit]
Description=Roon
After=network.target sound.target pipewire.service avahi-daemon.service
ConditionPathExists=/etc/uuid

[Service]
Type=simple
User=raat
Group=raat
Environment=HOME=/var/lib/raat
ExecStart=/usr/bin/start-raat.sh /var/lib/raat/hifiberry.conf
StandardOutput=journal
Restart=always
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
