# buildroot/package/hifiberry-mpd/mpd.service
[Unit]
Description=mpd
Wants=network.target
After=network.target sound.target

[Service]
Type=simple
ExecStartPre=/opt/hifiberry/bin/bootmsg "Starting music player daemon"
ExecStart=/bin/mpd --no-daemon

[Install]
WantedBy=multi-user.target
