# packages/hifiberry-tools/services/hifiberry-finish.service
[Unit]
Description=HiFiBerry final 
After=hifiberry.target
Before=psplash-quit.service

[Service]
Type=oneshot
ExecStartPre=/opt/hifiberry/bin/myip
ExecStartPre=sleep 1
ExecStart=/opt/hifiberry/bin/myurl
StandardOutput=journal

[Install]
WantedBy=multi-user.target
